# integrations/penguins-recovery/builders/lifeboat/zfiles/inittab
# /etc/inittab - Lifeboat Linux
# Auto-login root on tty1; additional consoles on tty2/tty3
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

tty1::respawn:/sbin/getty -n -l /bin/login 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3

::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown
