# systemd/ilf-update.service
[Unit]
Description=ILF atomic system upgrade
Documentation=man:ilf(1)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/ilf/ilf.toml

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ilf upgrade
StandardOutput=journal
StandardError=journal
# Prevent concurrent upgrade runs
ExecStartPre=/bin/sh -c '! systemctl is-active --quiet ilf-update.service'
