# integrations/penguins-recovery/builders/arch/airootfs/root/.config/autostart/easy-language-config.desktop
[Desktop Entry]
Type=Application
Name=Easy language config
Exec=/usr/local/bin/easy-language-config
StartupNotify=false
Terminal=false